# launchd/com.id7data.derive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.id7data.derive</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tomhausemer/Coding/ID7_data/.venv/bin/python</string>
        <string>/Users/tomhausemer/Coding/ID7_data/derive_sessions.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tomhausemer/Coding/ID7_data</string>

    <!-- Run every 900 seconds (15 minutes). -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tomhausemer/Coding/ID7_data/logs/derive.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tomhausemer/Coding/ID7_data/logs/derive.log</string>
</dict>
</plist>
